# systemd/ai-workflow-health.timer
[Unit]
Description=Workflow health check every 6 hours

[Timer]
OnCalendar=*-*-* 06/6:00:00
Persistent=true

[Install]
WantedBy=timers.target