# systemd/bazzite-mcp-bridge.service
[Unit]
Description=Bazzite MCP Bridge for Newelle
Documentation=file:///var/home/lch/projects/bazzite-laptop/docs/superpowers/specs/2026-03-25-newelle-mcp-bridge-router-overhaul-design.md
After=network.target

[Service]
Type=simple
WorkingDirectory=/var/home/lch/projects/bazzite-laptop
Environment=PATH=/var/home/lch/projects/bazzite-laptop/.venv/bin:/usr/bin:/bin
ExecStart=/var/home/lch/projects/bazzite-laptop/.venv/bin/python -m ai.mcp_bridge --bind 127.0.0.1 --port 8766
Restart=on-failure
RestartSec=5

# Security hardening
NoNewPrivileges=yes
ProtectSystem=strict
PrivateTmp=yes
ReadWritePaths=/var/home/lch/.config/bazzite-ai /var/home/lch/projects/bazzite-laptop /var/log/system-health /var/home/lch/security

[Install]
WantedBy=default.target
