# systemd/bazzite-llm-proxy.service
[Unit]
Description=Bazzite LLM Router Proxy for Newelle
Documentation=file:///var/home/lch/projects/bazzite-laptop/docs/superpowers/specs/2026-03-25-newelle-mcp-bridge-router-overhaul-design.md
After=network.target

[Service]
Type=simple
WorkingDirectory=/var/home/lch/projects/bazzite-laptop
Environment=PATH=/var/home/lch/projects/bazzite-laptop/.venv/bin:/usr/bin:/bin
ExecStart=/var/home/lch/projects/bazzite-laptop/.venv/bin/python -m ai.llm_proxy --bind 127.0.0.1 --port 8767
Restart=on-failure
RestartSec=10
RestartMax=5

Slice=ai.slice

# Resource limits — prevent memory leaks from freezing the system
MemoryAccounting=yes
MemoryHigh=80M
MemoryMax=150M
CPUWeight=80
TasksMax=128

# Security hardening
NoNewPrivileges=yes
ProtectSystem=strict
PrivateTmp=yes
ReadWritePaths=/var/home/lch/.config/bazzite-ai /var/home/lch/projects/bazzite-laptop

[Install]
WantedBy=default.target
