# systemd/clamav-healthcheck.service
[Unit]
Description=Bazzite Security Health Check

[Service]
Type=oneshot
ExecStart=/usr/local/bin/clamav-healthcheck.sh
User=root

# Security hardening
NoNewPrivileges=yes
ProtectSystem=strict
ProtectHome=read-only
PrivateTmp=yes
ReadWritePaths=/var/log/clamav-scans /home/lch/security
