# systemd/knowledge-storage.service
[Unit]
Description=Bazzite Knowledge Storage Health Check
After=network.target

[Service]
Type=oneshot
User=lch
WorkingDirectory=/home/lch/projects/bazzite-laptop
ExecStart=/usr/bin/python3 -m ai.agents.knowledge_storage
TimeoutStartSec=120
NoNewPrivileges=yes
PrivateTmp=yes
ReadWritePaths=/home/lch/security /var/log/system-health
StandardOutput=journal
StandardError=journal
SyslogIdentifier=knowledge-storage

[Install]
WantedBy=multi-user.target
