# systemd/lancedb-optimize.timer
[Unit]
Description=LanceDB weekly compaction timer

[Timer]
OnCalendar=Sun *-*-* 02:00:00
Persistent=true
RandomizedDelaySec=300

[Install]
WantedBy=timers.target