# systemd/rag-embed.service
# RAG nightly embedding service — chunk + embed security logs into LanceDB
# Deploy to: /etc/systemd/system/rag-embed.service
# Enable with: sudo systemctl enable rag-embed.timer

[Unit]
Description=RAG Security Intelligence — Nightly Embedding
After=network-online.target clamav-quick.service
Wants=network-online.target

[Service]
Type=oneshot
ExecStart=/usr/local/bin/rag-embed.sh --all
User=lch
Group=lch

# Low priority — never interfere with gaming or scans
Nice=19
IOSchedulingClass=idle
CPUSchedulingPolicy=idle

# Timeout: embedding should finish within 10 minutes
TimeoutStartSec=600

# Security hardening
NoNewPrivileges=yes
ProtectSystem=strict
ProtectHome=read-only
PrivateTmp=yes
ReadWritePaths=/home/lch/security/vector-db /home/lch/.config/bazzite-ai

# Logging
StandardOutput=journal
StandardError=journal
SyslogIdentifier=rag-embed
