# systemd/release-watch.timer
[Unit]
Description=Bazzite AI — daily upstream release watch (9:45 AM)

[Timer]
OnCalendar=*-*-* 09:45:00
Persistent=true

[Install]
WantedBy=timers.target
