# systemd/weekly-insights.service
[Unit]
Description=Weekly AI Insights Generator

[Service]
Type=oneshot
ExecStart=/home/lch/projects/bazzite-laptop/.venv/bin/python scripts/generate-weekly-insights.py
WorkingDirectory=/var/home/lch/projects/bazzite-laptop
Environment="PATH=/var/home/lch/projects/bazzite-laptop/.venv/bin:/usr/bin"
Environment="PYTHONPATH=/var/home/lch/projects/bazzite-laptop"

[Install]
WantedBy=multi-user.target
