# systemd/user/fedora-updates.service
[Unit]
Description=Bazzite AI — Fedora/Bodhi update awareness (User)
Documentation=file:///home/lch/projects/bazzite-laptop/docs/AGENT.md
After=network-online.target
Wants=network-online.target

[Service]
Type=oneshot
WorkingDirectory=%h/projects/bazzite-laptop
ExecStart=%h/projects/bazzite-laptop/.venv/bin/python -m ai.system.fedora_updates --check
StandardOutput=journal
StandardError=journal

# Conservative resource limits
Nice=10
IOSchedulingClass=best-effort
IOSchedulingPriority=7

[Install]
WantedBy=default.target
