# systemd/user/rag-embed.timer
[Unit]
Description=RAG Embedding — Daily 09:00 (User)

[Timer]
OnCalendar=*-*-* 09:00:00
Persistent=true
RandomizedDelaySec=300

[Install]
WantedBy=timers.target
